# voice-bridge.service
[Unit]
Description=Voice Bridge – Discord audio forwarding bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/voice-bridge-stable
ExecStart=/home/voice-bridge-stable/venv/bin/python /home/voice-bridge-stable/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
